# deploy/hollyness-api.service
# /etc/systemd/system/hollyness-api.service
# Install:
#   sudo cp /var/www/hollyness/deploy/hollyness-api.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now hollyness-api

[Unit]
Description=Hollyness & Respishers FastAPI Backend
After=network.target postgresql.service
Requires=postgresql.service

[Service]
Type=simple
User=root
WorkingDirectory=/var/www/hollyness/hollyness-api
EnvironmentFile=/var/www/hollyness/hollyness-api/.env

ExecStart=/var/www/hollyness/hollyness-api/venv/bin/uvicorn \
    app.main:app \
    --host 127.0.0.1 \
    --port 8000 \
    --workers 2 \
    --log-level info

Restart=on-failure
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=3

StandardOutput=journal
StandardError=journal
SyslogIdentifier=hollyness-api

NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
